# mix_job/config/deploy/templates/systemd/job.service.erb
[Unit]
Description=Background worker
After=network.target

[Service]
User=<%= fetch(:deployer_name) %>
WorkingDirectory=<%= fetch(:deploy_to) %>/current
ExecStart=/home/<%= fetch(:deployer_name) %>/.rbenv/bin/rbenv sudo RAKE_OUTPUT=true RAILS_ENV=<%= cap.env %> RAILS_APP=<%= cap.app %> bin/rake job:watch -- --queue=<%= fetch(:job_queue_name) %>
Restart=always
StandardInput=null
StandardOutput=file:<%= fetch(:deploy_to) %>/current/log/job.log
StandardError=file:<%= fetch(:deploy_to) %>/current/log/job.log
SyslogIdentifier=<%= fetch(:job_service) %>
KillMode=mixed
KillSignal=SIGTERM
SendSIGHUP=yes
FinalKillSignal=SIGINT
TimeoutStopSec=infinity

[Install]
WantedBy=multi-user.target
