# mr_recipe/config/provision/files/etc/systemd/system/postgresql.service.d/reload.conf
[Service]
ExecStartPost=/bin/systemctl reload postgresql
